# scripts/vps/nine-universe.service
[Unit]
Description=Nine Universe (Spring Boot)
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=nineuniverseapp
Group=nineuniverseapp
WorkingDirectory=/opt/nine-universe
Environment=SPRING_PROFILES_ACTIVE=production
# 1GB VPS 向けにヒープを抑えめに（余裕があれば 512m や 768m に）
ExecStart=/usr/bin/java -Xms128m -Xmx384m -jar /opt/nine-universe/nine-universe.jar
Restart=on-failure
RestartSec=10
EnvironmentFile=-/etc/nine-universe.env

[Install]
WantedBy=multi-user.target
